# fruit-default.openrc
#!/sbin/openrc-run

description="Deploy Fruit configurations."

depend() {
	after net root bootmisc modules syslog sysctl
}

start() {
	/usr/sbin/fruit-deploy
	return $?
}

stop() {
	return 0
}
