# openrc-fruit.in
#!/sbin/openrc-run

description="Deploy Fruit configurations."

depend() {
	before net
	after root bootmisc modules syslog sysctl
	keyword -lxc -docker
}

start() {
	/usr/sbin/fruit-deploy
	return $?
}

stop() {
	return 0
}
