# src/fruit-data
#!/sbin/openrc-run

share=/data/netboot/share

depend()
{
	before docker
	provide fruit-data
	need fruit-overlay
}

start()
{
	ebegin "Mounting nfs volume to /data, then setting up /var/lib/docker"
	srvip=$(mount | grep ' on / ' | cut -d' ' -f1 | cut -d':' -f1)
	ip=$(ip -4 addr show dev eth0 | grep inet | awk '{print $2}' | cut -d'/' -f1)
	[ -e /var/lib/docker ] && rm -rf /var/lib/docker
	mount -t nfs -o vers=4,rw,nolock ${srvip}:${share}/${ip} /data && \
	mkdir -p /data/docker && \
	ln -sf /data/docker /var/lib/docker
	ret=$?
	eend $ret

	return $ret
}

stop()
{
	umount -f /data
}
