# src/fruit-overlay
#!/sbin/openrc-run

description="Set up the overlayfs on /run, /var directories"

depend()
{
	need dev
	keyword -docker -lxc -prefix -systemd-nspawn -vserver
}

mount_tmpfs()
{
	ebegin "Mounting tmpfs on $1"
	mem=$(free -m | grep 'Mem: ' | awk '{print $2}')
	modprobe overlay && \
	mount -t tmpfs -o size=$(( $mem / 2 ))M tmpfs $1
	eend $?
}

mount_overlay()
{
	target="$1"
	overlaydir="$2"
	ebegin "Mounting overlay on $target using ramdisk at $2"
	upperdir=${overlaydir}/upperdir${target}
	lowerdir=${overlaydir}/lowerdir${target}
	mkdir -p ${upperdir} && \
	mkdir -p ${lowerdir} && \
	mount -t overlay -o lowerdir=${target},upperdir=${upperdir},workdir=${lowerdir} overlay ${target}
	eend $?
}

start()
{
	overlaydir="/media/overlay"
	mount_tmpfs ${overlaydir}
	mount_overlay "/var" ${overlaydir}
	mount_overlay "/etc" ${overlaydir}
	mount_overlay "/root" ${overlaydir}
	return 0
}
