# src/fruit-nat
#!/sbin/openrc-run

data="/run/fruit-nat"

depend()
{
  need net
  provide nat
}

install_iptables()
{
  [ "$(which iptables)" = "" ] && apk add iptables 1>/dev/null
}

netaddress() {
  addr=$(ip -4 addr show dev "$1" | grep 'inet ' | head -1 | awk '{print $2}')
  echo "$(ipcalc.sh ${addr} | cut -d' ' -f3)/$(echo ${addr} | cut -d'/' -f2)"
}

start()
{
  local wan=$(route | grep default | awk '{print $8}')
  local lan=$(cat /etc/network/interfaces | grep 'iface eth.* static' | awk '{print $2}')
  local ret=0

  ebegin "Enabling NAT lan=${lan} wan=${wan}"
  if [ "$wan" != "" ] && [ "$lan" != "" ] && [ "$wan" != "$lan" ]; then
    local lan_netaddr=$(netaddress "$lan")

    install_iptables
    sysctl -w net.ipv4.ip_forward=1 1>/dev/null && \
    iptables -A FORWARD -i ${lan} -p udp -m udp --dport 67 -j DROP 1>/dev/null 2>/dev/null && \
    iptables -A FORWARD -o ${wan} -p udp -m udp --dport 67 -j DROP 1>/dev/null 2>/dev/null && \
    iptables -A FORWARD -i ${wan} -o ${lan} -m state --state RELATED,ESTABLISHED -j ACCEPT 1>/dev/null 2>/dev/null && \
    iptables -A FORWARD -i ${lan} -o ${wan} -j ACCEPT 1>/dev/null 2>/dev/null && \
    iptables -t nat -A POSTROUTING -s ${lan_netaddr} -o ${wan} -j MASQUERADE 1>/dev/null 2>/dev/null
  
    ret=$?
    if [ $ret -eq 0 ]; then
      echo "$wan $lan $lan_netaddr" > "$data"
    fi
  else
    ret=1
  fi
  eend ${ret}
  return ${ret}
}

stop()
{
  local ret=0
  ebegin "Disabling NAT"
  if [ -f "$data" ]; then
    local wan=$(cat "$data" | cut -d' ' -f1)
    local lan=$(cat "$data" | cut -d' ' -f2)
    local lan_netaddr=$(cat "$data" | cut -d' ' -f3)

    install_iptables
    iptables -D FORWARD -i ${lan} -p udp -m udp --dport 67 -j DROP 1>/dev/null 2>/dev/null && \
    iptables -D FORWARD -o ${wan} -p udp -m udp --dport 67 -j DROP 1>/dev/null 2>/dev/null && \
    iptables -D FORWARD -i ${wan} -o ${lan} -m state --state RELATED,ESTABLISHED -j ACCEPT 1>/dev/null 2>/dev/null && \
    iptables -D FORWARD -i ${lan} -o ${wan} -j ACCEPT 1>/dev/null 2>/dev/null && \
    iptables -t nat -D POSTROUTING -s ${addr} -o ${wan} -j MASQUERADE 1>/dev/null 2>/dev/null && \

    ret=$?
  fi
  eend ${ret}
  return ${ret}
}
